# config/gmooc.conf
[program:gmooc]
command=/usr/local/bin/uwsgi --ini /root/project/Gmooc/uwsgi.ini
directory=/root/project/Gmooc
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
